# systemd/home-pc/roro-pull.timer
[Unit]
Description=Check for roro-bot updates every 60 seconds

[Timer]
OnBootSec=30
OnUnitActiveSec=60

[Install]
WantedBy=timers.target
