# systemd/home-pc/roro-pull.service
[Unit]
Description=Auto-pull roro-bot from GitHub

[Service]
Type=oneshot
ExecStart=/home/matt/.roro-bot/pull-and-restart.sh
